# tools/autostart/keyboard_hook_reader.desktop
[Desktop Entry]
Name=LinuxKeyboardHookReader
Type=Application
Exec=LinuxKeyboardHookReader
Terminal=true
